# scripts/launchd/com.mighty.claude-usage-token-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mighty.claude-usage-token-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mightydesigncenter/tools/claude-usage/scripts/token_refresh.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/mightydesigncenter/tools/claude-usage</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/claude-usage-token-refresh.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude-usage-token-refresh-err.log</string>
</dict>
</plist>
